# systemd/nox-daemon.service
[Unit]
Description=Nox Voice Daemon
After=network.target pulseaudio.service pipewire.service
Wants=pulseaudio.service

[Service]
Type=simple
ExecStart=%h/Documents/git/x11-mcp-voice/.venv/bin/python -m x11_mcp_voice
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:1

[Install]
WantedBy=default.target
